# packer/files/reddit.service
[Unit]
Description=reddit

[Service]
Type=simple
ExecStart=/home/ubuntu/start.sh

[Install]
WantedBy=multi-user.target
